# .config/systemd/user/bitlbee.socket
[Socket]
ListenStream=/run/user/1000/bitlbee.sock
Accept=yes
