# .config/systemd/user/firefoxcache.service
[Unit]
Description=mkdir firefox cache

[Service]
ExecStart=/usr/bin/mkdir /dev/shm/firefox_pim
Type=forking

[Install]
WantedBy=default.target
