# .config/systemd/user/gammastep-config.path
[Path]
PathModified=/home/pim/.config/gammastep/config.ini
