# .config/systemd/user/home-pim-Hosts-dawn.automount
# Genrated with Fscript

[Unit]
Description=Automount dawn
ConditionPathExists=/home/pim/Hosts/dawn

[Automount]
TimeoutIdleSec=180
Where=/home/pim/Hosts/dawn

[Install]
WantedBy=multi-user.target
