# .config/systemd/user/sway-config.path
[Path]
PathModified=/home/pim/.config/sway/config
